# deploy/systemd/disc-steward-watch.service
[Unit]
Description=Disc Steward watcher for completed rip folders
Wants=network-online.target
After=network-online.target
PartOf=disc-steward.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=ryan
Group=ryan
WorkingDirectory=/home/ryan/projects/disc-steward
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/ryan/projects/disc-steward/.venv/bin/disc-steward watch --config /home/ryan/projects/disc-steward/config.yaml --interval 30 --lock-file /run/disc-steward/disc-steward-watch.lock
Restart=always
RestartSec=15
RuntimeDirectory=disc-steward
RuntimeDirectoryMode=0755
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/home/ryan/projects/disc-steward /run/disc-steward

[Install]
WantedBy=disc-steward.target
